# etc/crontab
01 01 * * * /var/www/iopav21/bin/cronjob.daily.sh
